# root/etc/services.d/kanzi/run
#!/usr/bin/with-contenv bash
exec \
	s6-setuidgid abc gunicorn -b 0.0.0.0:8000 /app/kanzi/alexa:app